# deploy/ical-gunicorn.service
[Unit]
Description=ical gunicorn
After=network.target ical-env-fetch.service
Requires=ical-env-fetch.service

[Service]
User=django
Group=www-data
# New runtime files (cache, __pycache__, collectstatic) born 640/750, group
# www-data — no world-read, no per-deploy permission drift.
UMask=0027
WorkingDirectory=/var/www/django_websites/django-ical
# Secrets now come from SSM via the env-fetch oneshot (tmpfs), not an on-disk .env.
EnvironmentFile=/run/ical/.env
ExecStart=/var/www/django_websites/django-ical/.venv/bin/gunicorn django_myical.wsgi:application -c deploy/gunicorn.conf.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
